# update/selfhost/systemd/mos-update-agent.service
[Unit]
Description=My Own Suite self-host update agent
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=/etc/mos-update-agent.env
ExecStart=/usr/bin/env node __REPO_DIR__/update/selfhost/agent/mos-update-agent.cjs
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
